# systemd/archive_data.timer
[Unit]
Description=Periodically archive old data

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target